# dist/supervisord.ini
[supervisord]
nodaemon=true

[program:TODO_NOME_APP]
autorestart=true
command=TODO_COMANDO_PER_FARLO_PARTIRE
directory=TODO_DIRECTORY_DELLAPPLICAZIONE
stdout_logfile=syslog
stderr_logfile=syslog
user=TODO_UTENTE

[program:cron]
autorestart=true
command=crond -f
directory=/tmp
stdout_logfile=syslog
stderr_logfile=syslog
user=root

[program:rsyslog]
autorestart=true
command=rsyslogd -n
directory=/tmp
stdout_logfile=syslog
stderr_logfile=syslog
user=root
